# scripts/com.alta.stray_tripwire.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.alta.stray_tripwire</string>
    <!-- Catches the EXTERNAL writer that keeps creating stray files (findings.jsonl,
         messages_to_colin.json, research_queue.json, usage.json, "# GOALS.md") at the repo
         root instead of data/agent/. No writer in this repo is at fault (all paths are
         correctly anchored), so this watches the repo-root directory and fires the instant it
         changes; stray_tripwire.py records forensics (ctime, lsof, processes near ctime, content
         vs the data/agent canonical) to logs/stray_tripwire.log, then quarantines the stray.
         No-ops on unrelated root changes. `launchctl load` to enable. -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/taboost/quant/scripts/stray_tripwire.py</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/taboost/quant</string>
    <key>WatchPaths</key>
    <array><string>/Users/taboost/quant</string></array>
    <key>RunAtLoad</key><true/>
    <!-- WatchPaths only fires on direct-child changes of the repo root, so the
         tripwire could go weeks without running (last fire 2026-06-16). Hourly
         floor guarantees stray-job detection independent of file activity. -->
    <key>StartInterval</key><integer>3600</integer>
    <key>StandardOutPath</key><string>/Users/taboost/quant/logs/stray_tripwire.launchd.log</string>
    <key>StandardErrorPath</key><string>/Users/taboost/quant/logs/stray_tripwire.launchd.err</string>
    <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
