# scripts/com.alta.health.responder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.alta.health.responder</string>
    <!-- Component 1 of the autonomous layer. Every 30 min, 24/7: reads
         loop_health_status.json and writes a fix-request to dispatch_queue.jsonl
         for any DOWN loop. DISPATCH + NOTIFY ONLY — never restarts anything. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string><string>-c</string>
        <string>cd /Users/taboost/quant &amp;&amp; /opt/homebrew/bin/python3 sovereign/autonomous/health_responder.py</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/taboost/quant</string>
    <key>StartInterval</key><integer>1800</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/Users/taboost/quant/logs/health_responder.log</string>
    <key>StandardErrorPath</key><string>/Users/taboost/quant/logs/health_responder.err</string>
    <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
