# scripts/com.alta.decision_backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Entry-side decision-record backfill.
    16:45 ET Mon-Fri, after the session closes.

  WHY THIS EXISTS: scripts/backfill_decision_records.py creates minimal decision
  records for broker fills whose entry path never called log_forex_decision().
  It was added 2026-07-01 (28d1401) and NEVER SCHEDULED -- every record in
  data/decision_logs/ still carries "backfilled_at": "2026-07-01T06:11:46", the
  single manual run. Any unlogged fill since then has accumulated with no record
  at all, which is a silent hole in the label channel.

  Runs after com.alta.execution_harness (16:05) and before com.alta.system_eod
  (16:30 -- note: EOD reads labels, so this ideally moves earlier once verified;
  16:45 is deliberately AFTER for the first window so a bad backfill cannot
  corrupt an EOD note before anyone has reviewed its output).

  TRACKED-NOT-LOADED. Operator installs:
    cp scripts/com.alta.decision_backfill.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.alta.decision_backfill.plist
    python3 scripts/plist_watchdog.py --rebaseline "loaded decision_backfill"
-->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.alta.decision_backfill</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/taboost/quant/scripts/backfill_decision_records.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/taboost/quant</string>
  <key>EnvironmentVariables</key><dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONPATH</key><string>/Users/taboost/quant</string>
  </dict>
  <key>StandardOutPath</key><string>/Users/taboost/quant/logs/decision_backfill.log</string>
  <key>StandardErrorPath</key><string>/Users/taboost/quant/logs/decision_backfill.err</string>
  <key>Nice</key><integer>5</integer>
  <key>RunAtLoad</key><false/>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>45</integer><key>Weekday</key><integer>1</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>45</integer><key>Weekday</key><integer>2</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>45</integer><key>Weekday</key><integer>3</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>45</integer><key>Weekday</key><integer>4</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>45</integer><key>Weekday</key><integer>5</integer></dict>
  </array>
</dict>
</plist>
