# scripts/com.alta.weekly_review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.alta.weekly_review</string>
  <key>ProgramArguments</key><array><string>/opt/homebrew/bin/python3</string><string>/Users/taboost/quant/experience/weekly_review.py</string></array>
  <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
  <key>WorkingDirectory</key><string>/Users/taboost/quant</string>
  <key>Nice</key><integer>5</integer>
  <key>StandardOutPath</key><string>/Users/taboost/quant/logs/weekly_review.log</string>
  <key>StandardErrorPath</key><string>/Users/taboost/quant/logs/weekly_review.err</string>
  <key>StartCalendarInterval</key><dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
</dict>
</plist>
